# dev-lang/erlang/files/epmd.init-r2
#!/sbin/openrc-run
# Copyright 1999-2019 Gentoo Authors
# Distributed under the terms of the Erlang Public License 1.1

name="Erlang Port Mapper Daemon"

command="/usr/bin/epmd"
command_user="epmd"
command_group="epmd"

# epmd args, see: http://erlang.org/doc/man/epmd.html
#command_args="-address 127.0.0.1,::1"
# alternative: allow killing epmd even when there are active nodes:
command_args="-daemon -relaxed_command_check -address 127.0.0.1,::1"

command_background=yes
pidfile="/var/run/epmd.pid"

depend() {
	need net.lo
	before sshd
	after firewall
}

stop() {
	ebegin "Stopping ${name}"
	${command} -kill >/dev/null
	eend $?
}
