# dev-db/couchdb/files/couchdb.init-2
#!/sbin/openrc-run
# Copyright 1999-2019 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

# FIXME this file needs to be revised

pidfile=${COUCHDB_PID_FILE}
command=${EXEC:-/opt/couchdb/bin/couchdb}
command_args="-o ${COUCHDB_STDOUT_FILE} -e ${COUCHDB_STDERR_FILE} ${COUCHDB_OPTIONS}"
start_stop_daemon_args="--background --make-pidfile --user ${COUCHDB_USER}"

depend() {
	need net epmd
	use dns
	after firewall
}

start_pre() {
  checkpath -q -d -m 0755 -o ${COUCHDB_USER} $(dirname ${COUCHDB_PID_FILE})
}

stop_post() {
  if [ -f ${COUCHDB_PID_FILE} ]; then
    rm -f ${COUCHDB_PID_FILE}
  fi
}
